# debian/postinst
#!/bin/sh

set -e

if [ -d /usr/share/lider-cluster-server ]; then
	
	update-rc.d lider defaults
	systemctl enable lider.service
	systemctl start lider.service

fi

