# debian/postinst
#!/bin/sh

set -e

if [ -d /usr/share/lider-cluster-server ]; then
	systemctl enable lider.service
	systemctl start lider.service
fi

